# deploy/gymnasticon.sync-conflict-20260319-130814-FLRM4VV.service
[Unit]
Description=Gymnasticon Bike Bridge
After=bluetooth.target network.target
Wants=bluetooth.target

[Service]
Type=simple
Restart=always
RestartSec=5
Environment=NODE_ENV=production
WorkingDirectory=/opt/gymnasticon
# Mirror runtime logs to the console so HDMI output reflects what Gymnasticon emits during boot.
StandardOutput=journal+console
StandardError=journal+console
ExecStart=/usr/bin/node /opt/gymnasticon/src/app/cli.js --config /etc/gymnasticon.json

[Install]
WantedBy=multi-user.target
